# z.conf
[program:z]
command=python3 /home/curator2/z/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/z/test.err.log
stdout_logfile=/home/curator2/z/test.out.log